# services/disable-ati.service
[Unit]
Description=Disable ATI

[Service]
ExecStart=/home/ferhat/.X4/binaries/disable-ati

[Install]
WantedBy=multi-user.target
